# FASE-5/CÓDIGO_FASE_5/marketing-squad.service
[Unit]
Description=Marketing Squad · CMO + 11 AI Agents (Phase 5)
Documentation=https://github.com/Thiagaoai/agente-ads
After=network-online.target paperclip.service
Wants=network-online.target

[Service]
Type=simple
User=marketing
Group=marketing
WorkingDirectory=/opt/agente-ads/FASE-5/CÓDIGO_FASE_5
EnvironmentFile=/opt/agente-ads/FASE-5/CÓDIGO_FASE_5/.env
ExecStart=/usr/bin/node index.js

Restart=on-failure
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=300

StandardOutput=journal
StandardError=journal
SyslogIdentifier=marketing-squad

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/marketing-squad

LimitNOFILE=65536
MemoryMax=1G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
